# sysconf/profile.d/alias-fbterm-wallpaper.sh
alias fbterm='fbv -ciu /etc/splash/MyNewLinux/images/verbose-1280x800.jpg <<< q && FBTERM_BACKGROUND_IMAGE=1 fbterm'

